# resources/98-steelseries.rules
ACTION=="add", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1038" RUN+="/etc/udev/rules.d/98-steelseries-init.py '%E{DEVNAME}'"